# oc-patches/finalise/rc.local
# logread -f -F /board-resource/lastlog &
if [ -f /app/app ]; then
    echo 3 > /proc/sys/vm/drop_caches
    ulimit -c unlimited
    chmod +x /coredump_helper.sh
    echo "|/coredump_helper.sh %p %t" | tee /proc/sys/kernel/core_pattern
    chmod 755 /app/app
    /app/app&
fi

exit 0