# oc-patches/block-update-check-patch/rc.local
# Block update checks
iptables -A OUTPUT -d mms.chituiot.com -j REJECT
